# docker/supervisor.d/horizon.ini
[program:horizon]
command=php /site/artisan horizon
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_events_enabled=true
stderr_logfile_maxbytes=0
