# docker/supervisor.d/log.ini
[program:log]
command=/site/cli/docker link-logfile
user=nntmux
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_events_enabled=true
stderr_logfile_maxbytes=0
