# docker/supervisor.d/php-fpm.ini
[program:php-fpm]
command=php-fpm --nodaemonize -c /etc/php/php-fpm.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_events_enabled=true
stderr_logfile_maxbytes=0
